# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=info
pidfile=/tmp/supervisord.pid

[program:mongodb]
command=mongod --replSet rs0 --port %(ENV_MONGO_1_PORT)s --dbpath /data/db --bind_ip 127.0.0.1 --oplogSize 128
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:redis]
command=redis-server --port %(ENV_REDIS_PORT)s --dir /data/redis --appendonly yes --maxmemory %(ENV_REDIS_MAXMEMORY)s --maxmemory-policy noeviction --protected-mode no
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=3

[program:minio]
command=minio server /data/minio --console-address ":%(ENV_MINIO_WEB_PORT)s" --address ":%(ENV_MINIO_PORT)s"
environment=MINIO_ROOT_USER="%(ENV_AWS_ACCESS_KEY_ID)s",MINIO_ROOT_PASSWORD="%(ENV_AWS_SECRET_ACCESS_KEY)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=3

[program:create-bucket]
command=/bin/sh -c "until (mc alias set minio http://127.0.0.1:%(ENV_MINIO_PORT)s %(ENV_AWS_ACCESS_KEY_ID)s %(ENV_AWS_SECRET_ACCESS_KEY)s 2>/dev/null); do sleep 2; done && mc mb minio/%(ENV_MINIO_BUCKET)s 2>/dev/null; exit 0"
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:any-sync-coordinator]
command=any-sync-coordinator -c /data/etc/any-sync-coordinator/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:any-sync-node-1]
command=any-sync-node -c /data/etc/any-sync-node-1/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:any-sync-node-2]
command=any-sync-node -c /data/etc/any-sync-node-2/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:any-sync-node-3]
command=any-sync-node -c /data/etc/any-sync-node-3/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:any-sync-filenode]
command=any-sync-filenode -c /data/etc/any-sync-filenode/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:any-sync-consensusnode]
command=any-sync-consensusnode -c /data/etc/any-sync-consensusnode/config.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
